# raspi_pxe_docker_fleet/rootfs/usr/local/lib/ha-pxe/templates/ha-pxe-container-sync.service
[Unit]
Description=Reconcile Home Assistant managed Docker containers
Wants=network-online.target docker.service
After=network-online.target docker.service ha-pxe-firstboot.service
ConditionPathExists=/etc/ha-pxe/containers.json
ConditionPathExists=/var/lib/ha-pxe/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ha-pxe-container-sync
