# raspi_pxe_docker_fleet/rootfs/usr/local/lib/ha-pxe/templates/ha-pxe-firstboot.service
[Unit]
Description=Home Assistant PXE client first boot provisioning
Wants=network-online.target
After=network-online.target cloud-final.service
ConditionPathExists=/etc/ha-pxe/bootstrap.env
ConditionPathExists=!/var/lib/ha-pxe/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ha-pxe-firstboot
RemainAfterExit=yes
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
