# raspi_pxe_docker_fleet/rootfs/usr/local/lib/ha-pxe/templates/ha-pxe-container-sync.timer
[Unit]
Description=Recurring reconciliation for Home Assistant managed Docker containers

[Timer]
OnBootSec=45s
OnUnitActiveSec=15m
AccuracySec=1m
Persistent=true
Unit=ha-pxe-container-sync.service

[Install]
WantedBy=timers.target
