# raspi_pxe_docker_fleet/rootfs/usr/local/lib/ha-pxe/templates/ha-pxe-early-log.service
[Unit]
Description=Home Assistant PXE early boot diagnostics
Wants=network.target
After=network.target
Before=ha-pxe-firstboot.service
ConditionPathExists=/etc/ha-pxe/bootstrap.env
ConditionPathExists=!/var/lib/ha-pxe/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ha-pxe-early-log
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
